# templates/kroagent-dashboard.service
[Unit]
Description=KroAgent Dashboard
After=network.target

[Service]
Type=simple
User={{USER}}
Environment=DASHBOARD_PORT=18900
Environment=KROAGENTS_DIR={{INSTALL_DIR}}
Environment=DASHBOARD_BIND=127.0.0.1
Environment=PATH={{HOME}}/.npm-global/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 {{INSTALL_DIR}}/web/dashboard_server.py
Restart=on-failure
RestartSec=5
WorkingDirectory={{INSTALL_DIR}}

[Install]
WantedBy=multi-user.target
